# supervisord-apache2.conf
[program:apache2]
command=apache2 -DFOREGROUND
numprocs=1
autostart=true
autorestart=true
