# supervisord-otrscron.conf
[program:otrscron]
command=/otrscron.sh
numprocs=1
autostart=true
autorestart=true
