# supervisord-apache2.conf
[program:apache2]
command=apache2 -DFOREGROUND
stdout_events_enabled=true
stderr_events_enabled=true
numprocs=1
autostart=true
autorestart=true
